# scripts/ci-mac.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- launchd agent for the smedjan Mac GPU CI runner. Deploy on an always-on Apple-Silicon Mac:
       1. Clone a DEDICATED checkout:  git clone <origin> ~/ci/smedjan
       2. Substitute the home path:    sed "s#__HOME__#$HOME#g" ~/ci/smedjan/scripts/ci-mac.plist \
                                         > ~/Library/LaunchAgents/nu.smedjan.ci.plist
       3. Load it:                     launchctl load ~/Library/LaunchAgents/nu.smedjan.ci.plist
     It polls origin/main every 30 min and runs the full Mac CI gate only when it advances.
     Verdict + log: ~/.smedjan-ci-report.txt ; last-tested commit: ~/.smedjan-ci-last . -->
<plist version="1.0">
<dict>
  <key>Label</key><string>nu.smedjan.ci</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__HOME__/ci/smedjan/scripts/ci-mac-poll.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>__HOME__/.cargo/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>SMEDJAN_CI_REPO</key><string>__HOME__/ci/smedjan</string>
  </dict>
  <key>StartInterval</key><integer>1800</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>__HOME__/.smedjan-ci-launchd.log</string>
  <key>StandardErrorPath</key><string>__HOME__/.smedjan-ci-launchd.log</string>
</dict>
</plist>
